# usbip/rootfs/etc/services.d/usbip/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting usbip..."

exec socat PTY,raw,echo=0,link=/dev/$(bashio::config 'devId') tcp:$(bashio::config 'server')
