# usbip/rootfs/etc/services.d/usbip/run
#!/usr/bin/with-contenv bashio

devId=$(bashio::config 'devId')
server=$(bashio::config 'server')

bashio::log.info "Starting usbip..."

exec /usr/lib/linux-tools/$(uname -r)/usbip attach -r $server -b \
      $(/usr/lib/linux-tools/$(uname -r)/usbip list -r $server | grep "$devId" | cut -d: -f1)
